# in/unit/push-to-cloud.service.in
[Unit]
Description=All-of-Us Arrays Push-to-Cloud Daemon
After=network.target
StartLimitIntervalSec=0

[Service]
ExecStart=java -cp %PTC_BIN clojure.main -m ptc.start
Restart=always
RestartSec=1
SuccessExitStatus=143
User=ptc

[Install]
WantedBy=default.target
