# in/unit/push-to-cloud.service.in
[Unit]
Description=All-of-Us Arrays Push-to-Cloud Daemon
After=network.target ptc-vault.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
EnvironmentFile=-/etc/sysconfig/push-to-cloud
ExecStart=java -cp %PTC_BIN clojure.main -m ptc.start
Restart=always
RestartSec=60
SuccessExitStatus=143
User=ptc

[Install]
WantedBy=default.target
