# in/unit/push-to-cloud.service.in
[Unit]
Description=All-of-Us Arrays Push-to-Cloud Daemon

[Service]
ExecStart=%INSTALL_PREFIX%/bin/start

[Install]
WantedBy=default.target
